# templates/monocle_worker.ini
[program:%{WORKER_NAME}]
command=/tmp/python3/monocle/bin/python3 scan.py --signature=%{SERVER_PORT} --no-status-bar %{BOOTSTRAP} %{NO_PICKLE} --log-level INFO 
user=phonephyo
numprocs=1
directory=/var/www/Monocle-%{SERVER_PORT}/current
startsecs=5
startretries=0
autorestart=false
stopwaitsecs=30
stopasgroup=true
killasgroup=true
stdout_logfile=/var/www/Monocle-%{SERVER_PORT}/current/logs/worker.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=1
redirect_stderr=true

[program:%{WORKER_NAME}_web]
command=/tmp/python3/monocle/bin/python3 web.py -H 0.0.0.0 -P %{SERVER_PORT} 
user=phonephyo
numprocs=1
directory=/var/www/Monocle-%{SERVER_PORT}/current
startsecs=5
startretries=3
autorestart=true
stopwaitsecs=300
stdout_logfile=/var/www/Monocle-%{SERVER_PORT}/current/logs/web.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
redirect_stderr=true
