# templates/monocle_web.ini
[program:monocle_web]
command=/opt/python3-venv/bin/python3 web.py -H 0.0.0.0 -P 8888
user=phonephyo
numprocs=1
directory=/var/www/Monocle/
startsecs=5
startretries=3
autorestart=true
stopwaitsecs=300
stdout_logfile=/var/www/Monocle/logs/web.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=1
redirect_stderr=true
